# epm/dist/deb/debian/postrm
#!/bin/bash

EPM_CONFIG_DIR=/etc/epm
rm -f $EPM_CONFIG_DIR/config.toml
systemctl disable epm
systemctl stop epm